# roles/idp/files/cron.daily/updateIDPsecret
#!/bin/bash

/opt/shibboleth-idp/bin/updateIDPsecret.sh
